# etc/profile.d/silverblue-firstboot.sh
if test "$(id -u)" -gt "0" && test -d "$HOME"; then
    if test ! -e "$HOME"/.config/silverblue/firstboot-done; then
        mkdir -p "$HOME"/.config/autostart
        cp -f /etc/skel.d/.config/autostart/silverblue-firstboot.desktop "$HOME"/.config/autostart
    fi
fi
